# system/etc/cron.d/vegguide
30 10 * * * www-data  /opt/perl5.14.2-no-threads/bin/perl /usr/local/bin/vg-regen-data-feed-rss
15 12 * * 7 www-data  /opt/perl5.14.2-no-threads/bin/perl /usr/local/bin/vg-send-watch-list-emails
37 3  * * * autarch   /opt/perl5.14.2-no-threads/bin/perl /usr/local/bin/vg-backup
17 6  * * * www-data  /opt/perl5.14.2-no-threads/bin/perl /usr/local/bin/vg-update-from-external-sources
2  15 * * * autarch   /opt/perl5.14.2-no-threads/bin/perl /usr/local/bin/vg-delete-old-sessions
